# debian/live-guest@.service
[Unit]
After=%i.device

[Service]
Type=forking
Restart=no
ExecStart=/usr/share/live-guest/mount %I
